# services/utserver/run
#!/bin/bash

exec /sbin/setuser baca /app/utserver/utserver -settingspath /config/utserver