# services/utserver/run
#!/bin/bash

exec /sbin/setuser baca /app/utserver/utserver -settingspath /data -configfile /config/utserver.conf -logfile /config/utserver.log